# extra/wipewine.service
[Unit]
Description=Wine file associations wiper

[Service]
Type=notify
ExecStart=%h/bin/wipewine

[Install]
WantedBy=default.target
